# deploy/systemd/pointerverse-update.service
[Unit]
Description=Update Pointerverse from GitHub
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-pointerverse
StandardOutput=journal
StandardError=journal
